# services.d/mysqld/run
#!/bin/sh

mysqld --user=mysql --skip-name-resolve --bind-address=0.0.0.0 --pid-file=/var/run/mysqld/mysqld.pid --open-files-limit=5120
